# scripte/mocplay.service
[Unit]
Description=MOC Player
After=network.target sound.target

[Service]
Type=forking
User=hara
ExecStart=/usr/bin/mocp -S
ExecStop=/usr/bin/mocp -x
Restart=no
WorkingDirectory=/home/hara

[Install]
WantedBy=multi-user.target
